# files/systemd/media-bridge-fix-mac.service
[Unit]
Description=Fix bridge MAC address for DHCP persistence
After=systemd-networkd.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/media-bridge-fix-mac
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target