# scripts/systemd/ndi-display.service
[Unit]
Description=NDI Display Service - Receive and display NDI streams on HDMI outputs
After=network-online.target ndi-bridge.service
Wants=network-online.target
PartOf=ndi-bridge.service

[Service]
Type=simple
User=ndi-bridge
Group=video
Environment="HOME=/var/lib/ndi-bridge"
Environment="NDI_DISPLAY_CONFIG=/etc/ndi-bridge/display-config.json"

# Start with auto-mapping by default
ExecStart=/usr/local/bin/ndi-display auto

# Restart on failure
Restart=on-failure
RestartSec=5

# Resource limits
LimitNOFILE=65536
LimitRTPRIO=95
LimitMEMLOCK=infinity

# CPU scheduling for real-time video
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80

# Allow access to display devices
SupplementaryGroups=video render
DeviceAllow=/dev/dri/card* rw
DeviceAllow=/dev/fb* rw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndi-display

[Install]
WantedBy=multi-user.target