# files/systemd/system/ndi-capture.service
[Unit]
Description=NDI Capture Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mediabridge
Group=audio
Restart=always
RestartSec=5

Environment="LD_LIBRARY_PATH=/usr/local/lib"

ExecStart=/opt/media-bridge/run.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
