# files/systemd/system/ndi-display@.service
[Unit]
Description=NDI Display Output %i
After=network.target
Wants=ndi-display-monitor.service
Before=getty@tty%i.service
Conflicts=getty@tty%i.service

[Service]
Type=simple
Restart=on-failure
RestartSec=5
User=root

# Environment for display selection
Environment="DISPLAY_ID=%i"
Environment="LD_LIBRARY_PATH=/usr/local/lib"

# Read stream configuration from file
EnvironmentFile=-/etc/media-bridge/display-%i.conf

# Start NDI display using launcher script (handles spaces in stream names)
# Console management is now handled by the launcher itself
ExecStart=/usr/local/bin/ndi-display-launcher %i
ExecStopPost=/usr/bin/rm -f /var/run/ndi-display/display-%i.status
ExecStopPost=/bin/sh -c 'echo 1 > /sys/class/vtconsole/vtcon1/bind 2>/dev/null || true'
ExecStopPost=/bin/sh -c 'chvt 1; sleep 0.1; chvt 2'
ExecStopPost=/bin/sh -c 'setterm -blank poke > /dev/tty2 2>/dev/null || true'

# Resource limits
LimitNOFILE=65536
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndi-display-%i

[Install]
WantedBy=multi-user.target