# files/systemd/system/ndi-bridge-collector.service
[Unit]
Description=NDI Bridge Metrics Collector
After=ndi-bridge.service
Wants=ndi-bridge.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ndi-bridge-collector
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target