# scripts/helper-scripts/media-bridge-intercom.service
[Unit]
Description=Media Bridge Intercom Service
After=network-online.target pipewire-system.service pipewire-pulse-system.service wireplumber-system.service
Wants=network-online.target
Requires=pipewire-system.service pipewire-pulse-system.service wireplumber-system.service

[Service]
Type=simple
User=root

# Pre-start cleanup
ExecStartPre=-/bin/bash -c 'pkill -9 -f chrome 2>/dev/null || true'
ExecStartPre=-/bin/bash -c 'pkill -9 -f x11vnc 2>/dev/null || true'
ExecStartPre=-/bin/bash -c 'pkill -9 -f Xvfb 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'rm -f /tmp/.X99-lock /tmp/chrome-vdo-profile/Singleton* 2>/dev/null || true'
ExecStartPre=/bin/sleep 3

# Setup virtual audio devices
ExecStartPre=/usr/local/bin/media-bridge-audio-manager setup

# Start the service (auto-selects PipeWire or PulseAudio)
ExecStart=/usr/local/bin/media-bridge-intercom-launcher

# Stop cleanup
ExecStopPost=-/bin/bash -c 'pkill -9 -f chrome 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f x11vnc 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f Xvfb 2>/dev/null || true'

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target