# scripts/helper-scripts/media-bridge-intercom.service
[Unit]
Description=Media Bridge Intercom Service
After=network-online.target pipewire-system.service pipewire-pulse-system.service
Wants=network-online.target
Requires=pipewire-system.service pipewire-pulse-system.service

[Service]
Type=simple
TimeoutStartSec=300
# Run as mediabridge user with all services unified
User=mediabridge
Group=audio
SupplementaryGroups=video pipewire

# Environment variables for mediabridge user
Environment="HOME=/var/lib/mediabridge"
Environment="XDG_RUNTIME_DIR=/run/user/999"
Environment="PULSE_RUNTIME_PATH=/run/user/999/pulse"
Environment="PULSE_SERVER=unix:/run/user/999/pulse/native"

# Setup virtual audio devices (run as root)
ExecStartPre=+/usr/local/bin/media-bridge-audio-manager setup

# Start the service as intercom user
ExecStart=/usr/local/bin/media-bridge-intercom-launcher

# Stop cleanup
ExecStopPost=-/bin/bash -c 'pkill -9 -f chrome 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f x11vnc 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f Xvfb 2>/dev/null || true'

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target