# scripts/helper-scripts/media-bridge-intercom.service
[Unit]
Description=Media Bridge Intercom Service
After=network-online.target
Wants=network-online.target
# PipeWire now runs as user service, no direct dependencies

[Service]
Type=simple
User=mediabridge
Group=audio

# Environment for PipeWire socket access
Environment="XDG_RUNTIME_DIR=/run/pipewire"
Environment="PIPEWIRE_RUNTIME_DIR=/run/pipewire"
Environment="PULSE_RUNTIME_PATH=/run/pipewire/pulse"

# Chrome profile location
Environment="CHROME_USER_DATA_DIR=/var/lib/mediabridge/chrome-profile"

# Pre-start cleanup
ExecStartPre=-/bin/bash -c 'pkill -9 -f chrome 2>/dev/null || true'
ExecStartPre=-/bin/bash -c 'pkill -9 -f x11vnc 2>/dev/null || true'
ExecStartPre=-/bin/bash -c 'pkill -9 -f Xvfb 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'rm -f /tmp/.X99-lock /var/lib/mediabridge/chrome-profile/Singleton* 2>/dev/null || true'
ExecStartPre=/bin/sleep 3

# Setup virtual audio devices (now runs as mediabridge)
ExecStartPre=/usr/local/bin/media-bridge-audio-manager setup

# Start the service
ExecStart=/usr/local/bin/media-bridge-intercom-launcher

# Stop cleanup
ExecStopPost=-/bin/bash -c 'pkill -9 -f chrome 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f x11vnc 2>/dev/null || true'
ExecStopPost=-/bin/bash -c 'pkill -9 -f Xvfb 2>/dev/null || true'

# Restart policy
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target