# scripts/helper-scripts/ndi-bridge-intercom.service
[Unit]
Description=NDI Bridge VDO.Ninja Intercom Service
After=network-online.target ndi-bridge.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
ExecStart=/opt/ndi-bridge/ndi-bridge-intercom
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndi-bridge-intercom

# Resource limits
MemoryLimit=512M
CPUQuota=20%

# Environment
Environment="DISPLAY=:0"
Environment="HOME=/root"

[Install]
WantedBy=multi-user.target