# scripts/helper-scripts/pipewire-system.service
[Unit]
Description=PipeWire System Service
After=sound.target pipewire-system.socket systemd-tmpfiles-setup.service
Before=media-bridge-intercom.service ndi-display@0.service ndi-display@1.service ndi-display@2.service
Requires=pipewire-system.socket
Wants=systemd-tmpfiles-setup.service

[Service]
Type=simple
User=root
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PIPEWIRE_RUNTIME_DIR=/run/user/0"
Environment="XDG_STATE_HOME=/var/lib/pipewire"
Environment="XDG_CACHE_HOME=/var/cache/pipewire"
Environment="HOME=/root"
ExecStartPre=/bin/mkdir -p /run/user/0
ExecStartPre=/bin/mkdir -p /run/user/0/pulse
ExecStartPre=/bin/mkdir -p /var/lib/pipewire
ExecStartPre=/bin/mkdir -p /var/cache/pipewire
# Small delay to ensure system is ready
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/pipewire -c /usr/share/pipewire/pipewire.conf
# Force PipeWire restart if socket doesn't exist
ExecStartPost=/bin/bash -c 'sleep 5; if [ ! -S /run/user/0/pipewire-0 ]; then systemctl restart pipewire-system; fi'
Restart=always
RestartSec=5
# High priority for audio processing
Nice=-11
LimitRTPRIO=95
LimitNICE=-11
LimitMEMLOCK=256M

[Install]
WantedBy=multi-user.target