# scripts/helper-scripts/pipewire-system.service
[Unit]
Description=PipeWire System Service
After=sound.target pipewire-system.socket systemd-tmpfiles-setup.service user-runtime-dir@0.service
Before=media-bridge-intercom.service ndi-display@0.service ndi-display@1.service ndi-display@2.service
Requires=pipewire-system.socket user-runtime-dir@0.service
Wants=systemd-tmpfiles-setup.service

[Service]
Type=simple
User=root
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PIPEWIRE_RUNTIME_DIR=/run/user/0"
Environment="XDG_STATE_HOME=/var/lib/pipewire"
Environment="XDG_CACHE_HOME=/var/cache/pipewire"
Environment="HOME=/root"
ExecStartPre=/bin/mkdir -p /run/user/0
ExecStartPre=/bin/mkdir -p /run/user/0/pulse
ExecStartPre=/bin/mkdir -p /var/lib/pipewire
ExecStartPre=/bin/mkdir -p /var/cache/pipewire
# Create runtime directory and wait for system to stabilize
ExecStartPre=/bin/bash -c 'mkdir -p /run/user/0 && chmod 700 /run/user/0 && sleep 5'
ExecStart=/usr/bin/pipewire -c /etc/pipewire/pipewire-system.conf
# Log socket status for debugging
ExecStartPost=/bin/bash -c 'sleep 3; ls -la /run/user/0/pipewire* 2>&1 | logger -t pipewire-socket-check'
Restart=always
RestartSec=5
# High priority for audio processing
Nice=-11
LimitRTPRIO=95
LimitNICE=-11
LimitMEMLOCK=256M

[Install]
WantedBy=multi-user.target