# scripts/helper-scripts/pipewire-pulse-system.service
[Unit]
Description=PipeWire PulseAudio System Service
After=pipewire-system.service
Requires=pipewire-system.service
BindsTo=pipewire-system.service

[Service]
Type=simple
User=mediabridge
Group=audio
SupplementaryGroups=video pipewire
Environment="XDG_RUNTIME_DIR=/run/user/999"
Environment="PULSE_RUNTIME_PATH=/run/user/999/pulse"
Environment="XDG_CACHE_HOME=/var/lib/mediabridge/.cache/pipewire-pulse"
Environment="HOME=/var/lib/mediabridge"
Environment="PIPEWIRE_RUNTIME_DIR=/run/user/999"
ExecStartPre=/bin/mkdir -p /run/user/999/pulse
ExecStartPre=/bin/mkdir -p /var/lib/mediabridge/.cache/pipewire-pulse
ExecStart=/usr/bin/pipewire-pulse
Restart=always
RestartSec=5
# High priority for audio processing
Nice=-11
LimitRTPRIO=95
LimitNICE=-11

[Install]
WantedBy=multi-user.target