# scripts/helper-scripts/vdo-ninja-intercom-usermode.service
[Unit]
Description=VDO.Ninja Intercom Service (User Mode)
After=network.target sound.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'pkill -u intercom 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'rm -f /tmp/.X99-lock 2>/dev/null || true'
ExecStart=/usr/local/bin/vdo-ninja-intercom-usermode
Restart=always
RestartSec=10
StartLimitInterval=0
StartLimitBurst=0

# Run as root to manage the intercom user
User=root
Group=root

# Allow access to audio devices
SupplementaryGroups=audio

# Environment
Environment="HOME=/root"
Environment="VDO_ROOM=nl_interkom"

[Install]
WantedBy=multi-user.target