# web/ndi-bridge-intercom-web.service
[Unit]
Description=NDI Bridge Intercom Web Interface
After=network.target ndi-bridge-intercom.service
Wants=ndi-bridge-intercom.service

[Service]
Type=exec
User=root
WorkingDirectory=/opt/ndi-bridge-web/backend
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PULSE_RUNTIME_PATH=/run/user/0/pulse"

[Install]
WantedBy=multi-user.target